# deploy/azure-vm/group-stock-backend.service
[Unit]
Description=Group Stock Django Backend
After=network.target

[Service]
User=azureuser
Group=www-data
WorkingDirectory=/opt/group_stock_project/backend
EnvironmentFile=/etc/group-stock/backend.env
Environment=APP_DIR=/opt/group_stock_project
Environment=VENV_DIR=/opt/group_stock_project/.venv
ExecStart=/opt/group_stock_project/deploy/azure-vm/gunicorn_start.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
